# raspberry-files/kiosk.service
[Unit]
Description=Kiosk Video Player
After=network.target

[Service]
Type=simple
User=infoactive
Group=infoactive
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/infoactive/.Xauthority
Environment=HOME=/home/infoactive
ExecStart=/bin/bash /var/www/kiosk/player/start-kiosk.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
