# raspberry-files/kiosk.service
[Unit]
Description=Kiosk Video Player
Wants=network-online.target
After=network-online.target
After=graphical.target
After=multi-user.target
After=lightdm.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=HOME=/home/pi
WorkingDirectory=/var/www/kiosk
ExecStartPre=/bin/mkdir -p /var/www/kiosk
ExecStartPre=/bin/chown -R pi:pi /var/www/kiosk
ExecStartPre=/bin/touch /var/log/kiosk.log
ExecStartPre=/bin/chown pi:pi /var/log/kiosk.log
ExecStartPre=/bin/sleep 10
ExecStart=/bin/bash /var/www/kiosk/start-kiosk.sh
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target
