# raspberry-files/kiosk.service
[Unit]
Description=Kiosk Video Player
After=network.target graphical.target

[Service]
Type=simple
User=infoactive
Group=infoactive
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/infoactive/.Xauthority
Environment=HOME=/home/infoactive
WorkingDirectory=/var/www/kiosk

# Start the kiosk
ExecStart=/bin/bash /var/www/kiosk/start-kiosk.sh

# Restart settings
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target