# raspberry-files/kiosk.service
[Unit]
Description=Kiosk Video Player
After=network.target device-monitor.service

[Service]
Type=simple
User=infoactive
Group=infoactive
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/infoactive/.Xauthority
Environment=HOME=/home/infoactive
WorkingDirectory=/var/www/kiosk

# Start the kiosk browser
ExecStart=/usr/bin/chromium-browser --noerrdialogs --disable-infobars --kiosk http://localhost/index.html

# Restart settings
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=3

[Install]
WantedBy=multi-user.target