# raspberry-files/monitor.service
[Unit]
Description=Kiosk Device Monitor and Browser Manager
After=network.target

[Service]
Type=simple
User=pi
Environment=BACKEND_URL=https://vinculo.com.py/new-player/api
ExecStart=/bin/bash -c '/usr/bin/python3 /var/www/kiosk/device_monitor.py & /usr/bin/python3 /var/www/kiosk/manage_browser.py'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
